# deploy/ec2/threadcast.service
[Unit]
Description=ThreadCast Spring Boot Application
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/threadcast
ExecStart=/usr/bin/java -jar -Xms256m -Xmx512m -Dspring.profiles.active=prod /home/ec2-user/threadcast/threadcast-server.jar
SuccessExitStatus=143
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
EnvironmentFile=/home/ec2-user/threadcast/.env

[Install]
WantedBy=multi-user.target
